# root/etc/services.d/knxd/run
#!/usr/bin/with-contenv bash
exec \
	s6-setuidgid abc /usr/local/bin/knxd /config/knxd.ini
